# ansible/roles/monitoring/templates/node-exporter-units.timer.j2
[Unit]
Description=Periodically export systemd unit states for node_exporter

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
