# ansible/roles/backup/templates/restic-backup.timer.j2
[Unit]
Description=Daily off-site restic backup

[Timer]
OnCalendar={{ backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
